# roles/prometheus/templates/prometheus.conf.j2
start on startup

script
   /usr/bin/prometheus
end script
